# scripts/metrics/node_exporter.service
[Unit]
Description=Prometheus Node Exporter

[Service]
ExecStart=/usr/bin/node_exporter 
Restart=always

[Install]
WantedBy=multi-user.target
